# recipes-upstreamable/alljoyn/alljoyn/alljoyn.upstart
# Copyright (c) 2013 LG Electronics, Inc.
description "AllJoyn Daemon"

start on tv-ready

console output

respawn

script
	/usr/sbin/alljoyn-daemon --internal --no-bt --no-ice
end script
